# recipes/setup-apps-services/files/rugix-app-gc.service
[Unit]
Description=Garbage collect old Rugix app generations

[Service]
Type=oneshot
ExecStart=rugix-ctrl apps gc
